# server/monitor/crontab
* * * * * /app/monitor >> /var/log/cron.log 2>&1